# etc/system/systemd/fstrim.timer
[Unit]
Description=run fstrim.service once in a month

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
