# systemd_user_units/console.target
[Unit]
Description=Non-graphical services
Requires=dbus.socket
After=dbus.socket

[Install]
Alias=default.target

