# systemd_user_units/env.service
[Unit]
Description=User environment variables
Before=default.target

[Service]
Type=oneshot
ExecStart=-/usr/bin/systemctl --user set-environment DBUS_SESSION_BUS_ADDRESS=%t/dbus/user.socket SSH_AUTH_SOCK=%t/ssh_auth_sock

[Install]
WantedBy=default.target

