# systemd_user_units/sockets.target
[Unit]
Description=Sockets
Documentation=man:systemd.special(7)


