# systemd_user_units/ssh-agent.socket
[Unit]
Description=ssh-agent socket activation


[Socket]
ListenStream=%t/urxvtd-%H

[Install]
WantedBy=sockets.target

