# home/systemd_user_units/console.target
[Unit]
Description=Non-graphical services
Requires=dbus.service
Before=wm.target
After=dbus.service

[Install]
Alias=default.target

