# home/systemd_user_units/kalu.service
[Unit]
Description=Upgrade notifier w/ AUR support
Documentation=man:kalu(1)

[Service]
ExecStart=/usr/bin/kalu 

[Install]
WantedBy=wm.target




