# home/systemd_user_units/mate-settings-daemon.service
[Unit]
Description=Mate settings daemon

[Service]
ExecStart=/usr/lib/mate-settings-daemon/mate-settings-daemon 

[Install]
WantedBy=wm.target
