# home/systemd_user_units/pasystray.service
[Unit]
Description=Pulseaudio system tray
Documentation=man:pasystray(1)

[Service]
ExecStart=/usr/bin/pasystray

[Install]
WantedBy=wm.target




