# home/systemd_user_units/pasystray.service
[Unit]
Description=Pulseaudio system tray

[Service]
ExecStart=/usr/bin/pasystray

[Install]
WantedBy=wm.target




