# home/systemd_user_units/psd-resync.timer
[Unit]
Description=Timer for Profile-sync-daemon - 1Hour
PartOf=psd-resync.service psd.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=20min



