# home/systemd_user_units/ssh-add.service
[Unit]
Description=ssh-add

[Service]
ExecStart=/usr/bin/ssh-add ~/.config/ssh/id_rsa 


[Install]
WantedBy=wm.target




