# home/systemd_user_units/wm.target
[Unit]
Description=Window manager
After=console.target 



