# home/systemd_user_units/wm.target
[Unit]
Description=Window manager target
After=console.target 



