# etc/user@.service.d/devmon.conf
[Service]
EnvironmentFile=%h/.config/devmon/devmon.conf
