# etc/user@.service.d/gpg-agent.conf
[Service]
Environment=GNUPGHOME=%h/.config/gnupg
