# etc/user@.service.d/gpg-agent.conf
[Service]
Environment=GPG_AGENT_INFO=%h/.config/gnupg/gpg-agent.info
