# etc/user@.service.d/ssh-agent.conf
[Service]
Environment=SSH_AUTH_SOCK=/run/user/%I/ssh_auth_sock
