# etc/user@.service.d/urxvtd.conf
[Service]
Environment=RXVT_SOCKET=/run/user/%I/urxvtd-hortensia
